# resources/com.user.mac-file-organizer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.mac-file-organizer</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/mac-file-organizer</string>
        <string>--daemon</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>~/Library/Logs/mac-file-organizer-error.log</string>
    <key>StandardOutPath</key>
    <string>~/Library/Logs/mac-file-organizer-output.log</string>
</dict>
</plist>